# etc/rust-socket-activation.service
[Service]
# Below is replaced in Makefile
ExecStart=EXECSTART